# docker/supervisor.conf
[program:slask]
command = python /home/docker/code/slask.py
